# patient_info.service
[Unit]
Description=Gunicorn instance for your_app_name
After=network.target

[Service]
User=your_app_user
WorkingDirectory=/home/your_app_user/your_app_name
EnvironmentFile=/home/your_app_user/your_app_name/.env
ExecStart=/home/your_app_user/your_app_name/venv/bin/gunicorn --workers 3 --bind unix:/run/your_app_name.sock app:app
Restart=on-failure

[Install]
WantedBy=multi-user.target